# samplepi.service
[Unit]
Description=SamplePi Audio Playback System
After=multi-user.target

[Service]
Type=simple
User=%USER%
WorkingDirectory=%HOME%/SamplePi
# Uncomment these lines once framebuffer display is configured:
# Environment="SDL_FBDEV=/dev/fb0"
# Environment="SDL_VIDEODRIVER=fbcon"
ExecStart=%HOME%/SamplePi/.venv/bin/python3 -m samplepi.main
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
