# samplepi.service
[Unit]
Description=SamplePi Audio Playback System
After=graphical.target
Wants=graphical.target

[Service]
Type=simple
User=%USER%
WorkingDirectory=%HOME%/SamplePi
# Set DISPLAY for desktop mode (X11)
Environment="DISPLAY=:0"
Environment="XAUTHORITY=/home/%USER%/.Xauthority"
ExecStart=%HOME%/SamplePi/.venv/bin/python3 -m samplepi.main
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
